# other/bin/ttycol.service
[Unit]
Description=Change tty colours
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ttycol.sh

[Install]
WantedBy=multi-user.target
